# files/restic/restic.timer
[Unit]
Description=Run backup daily, defined by luouelle

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
