# units/nextcloud/nextcloud_backup.timer
[Unit]
Description=Run nextcloud app and postgres db backup hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
