# restic/dot-config/containers/systemd/restic.service
[Unit]
Description=Container for NAS backup script, wrapper function
After=check@prune
