# units/torrents/deluge_state_backup.service
[Unit]
Description=Backup deluge state for torrents user

[Service]
Type=oneshot
ExecStart=%h/deluge_state_backup.sh
