# units/torrents/deluge_state_backup.service
[Unit]
Description=Backup deluge state

[Service]
Type=oneshot
ExecStart=%h/backup.sh
