# units/torrents/deluge_state_backup.timer
[Unit]
Description=Run deluge state backup hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
