# units/torrents/media_backup.service
[Unit]
Description=Backup data dir for torrents user

[Service]
Type=oneshot
ExecStart=%h/media_backup.sh
